# AFIRGEN FINAL/main backend/afirgen.service
[Unit]
Description=AFIRGen Backend Service - Automated FIR Generation
Documentation=https://github.com/your-repo/afirgen
After=network.target

[Service]
Type=simple
User=ssm-user
WorkingDirectory=/opt/afirgen-backend/AFIRGEN FINAL/main backend
EnvironmentFile=/opt/afirgen-backend/AFIRGEN FINAL/main backend/.env
ExecStart=/opt/afirgen-backend/AFIRGEN FINAL/main backend/venv/bin/uvicorn agentv5:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
